# profile.d/PythonCCU-addon.sh
export PATH=${PATH}:/usr/local/addons/PythonCCU/bin
export PYTHONUSERBASE=/usr/local/addons/PythonCCU/lib/python3.9/site-packages
export LD_LIBRARY_PATH=${LD_LIBRARY_PATH}:/usr/local/addons/PythonCCU/lib
